# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:spora-web]
command=frankenphp run --config /app/frankenphp.conf --adapter caddyfile
autostart=true
autorestart=true
startsecs=2
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:spora-worker]
command=php /app/bin/spora worker:run --daemon
autostart=true
autorestart=true
startsecs=3
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
